# etc/udev-rules/15-huawei-e1750.rules
SUBSYSTEM=="usb", ATTRS{idProduct}=="1446", ATTRS{idVendor}=="12d1", RUN+="/usr/sbin/usb_modeswitch"
