# etc/udev-rules/20-huawei-e1750-rename.rules
ACTION=="add", KERNEL=="ttyUSB*", RUN+="/usr/local/bin/modem-symlink-add.rb %k"
