# etc/udev-rules/20-huawei-e1750-remove.rules
ACTION=="remove", KERNEL=="ttyUSB*", RUN+="/usr/local/bin/modem-symlink-remove.sh %k"
